# tools/frr.service.in
[Unit]
Description=FRRouting
Documentation=https://frrouting.readthedocs.io/en/latest/setup.html
After=networking.service systemd-sysctl.service csmgrd.service
Before=lttng-sessiond.service
OnFailure=heartbeat-failed@%n
Before=routing_telemetry.service
Wants=routing_telemetry.service

[Service]
Nice=-5
Type=forking
NotifyAccess=all
StartLimitInterval=3m
StartLimitBurst=3
TimeoutSec=@TIMEOUT_MIN@m
WatchdogSec=60s
RestartSec=5
Restart=always
LimitNOFILE=2048
LimitCORE=536870912
ExecStart=@e_sbindir@/frrinit.sh start
ExecStop=@e_sbindir@/frrinit.sh stop
ExecReload=@e_sbindir@/frrinit.sh reload

[Install]
WantedBy=multi-user.target
